# packages/kernel/fcpci/files/etc/hotplug.d/.svn/all-wcprops
K 25
svn:wc:ra_dav:version-url
V 65
/openwrt/!svn/ver/12005/trunk/package/madwifi/files/etc/hotplug.d
END
