# packages/kernel/fcpci/files/etc/hotplug.d/.svn/entries
9

dir
14679
https://svn.openwrt.org/openwrt/trunk/package/madwifi/files/etc/hotplug.d
https://svn.openwrt.org/openwrt



2008-07-30T05:22:16.679739Z
12005
nbd


svn:special svn:externals svn:needs-lock











3c298f89-4303-0410-b956-a3cf2f4a3e73

net
dir

