# packages/utils/luci-app-ffwizard-pberg/uci-defaults/device-defaults
#!/bin/sh

(
while true ; do
sleep 30
wifie=0

set_default_config()
{
	local cfg="$1"
	#local type=""
	config_get type "$cfg" "type"
	echo $type
	case $type in
		atheros)
			uci set freifunk.wifi_device.channel=1
			uci set freifunk.wifi_device.diversity=1
			uci set freifunk.wifi_device.disabled=0
			uci set freifunk.wifi_device.txpower=""
			uci set freifunk.wifi_device.country=276
			uci set freifunk.wifi_device.regdomain="0x37"
			uci set freifunk.wifi_iface=defaults
			uci set freifunk.wifi_iface.mode=adhoc
			uci set freifunk.wifi_iface.bssid="02:CA:FF:EE:BA:BE"
			uci set freifunk.wifi_iface.sw_merge=1
			uci commit freifunk
		;;
		mac80211)
			uci set freifunk.wifi_device.channel=1
			uci set freifunk.wifi_device.diversity=""
			uci set freifunk.wifi_device.disabled=0
			uci set freifunk.wifi_device.txpower=""
			uci set freifunk.wifi_device.country=DE
			uci set freifunk.wifi_iface=defaults
			uci set freifunk.wifi_iface.mode=adhoc
			uci set freifunk.wifi_iface.bssid="02:CA:FF:EE:BA:BE"
			uci set freifunk.wifi_iface.sw_merge=""
			uci commit freifunk
		;;
		broadcom)
			uci set freifunk.wifi_device.channel=1
			uci set freifunk.wifi_device.diversity=""
			uci set freifunk.wifi_device.disabled=0
			uci set freifunk.wifi_device.txpower=""
			uci set freifunk.wifi_device.country=DE
			uci set freifunk.wifi_device.txantenna=0
			uci set freifunk.wifi_device.rxantenna=0
			uci set freifunk.wifi_iface=defaults
			uci set freifunk.wifi_iface.mode=adhoc
			uci set freifunk.wifi_iface.bssid="02:CA:FF:EE:BA:BE"
			uci set freifunk.wifi_iface.sw_merge=""
			uci commit freifunk
		;;
	esac
}

config_load wireless && wifie=1

if [ $wifie -eq 1 ] ; then
config_foreach set_default_config wifi-device

uci set freifunk.wizard.latitude="52.47598"
uci set freifunk.wizard.longitude="13.40168"
uci commit freifunk
exit 0
fi
done
) >/dev/null &
